# infra/systemd/jovio-onboarding.service
# Jovio Onboarding Emails — runs once daily via the matching .timer
#
# Install:
#   sudo cp infra/systemd/jovio-onboarding.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now jovio-onboarding.timer
#   # Verify:
#   sudo systemctl list-timers | grep jovio
#   # Manual run:
#   sudo systemctl start jovio-onboarding.service

[Unit]
Description=Jovio onboarding email sequence (day 0 / 3 / 10 / 14)
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/jovi/api-server
EnvironmentFile=/etc/jovio/api-server.env
ExecStart=/home/ubuntu/jovi/api-server/node_modules/.bin/ts-node src/jobs/onboarding-emails.ts

StandardOutput=append:/var/log/jovio-onboarding.log
StandardError=append:/var/log/jovio-onboarding.log

# Failure email to root if it errors — CloudWatch will alert via the log group
OnFailure=

[Install]
WantedBy=multi-user.target
