# infra/systemd/jovio-api.service
# Jovio API Server (Node + Express + TypeScript)
#
# Install:
#   sudo cp infra/systemd/jovio-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now jovio-api
#
# NOTE: previous Supervisor config ran `node src/index.ts` which DOES NOT WORK
# because Node can't run TS directly. The correct ExecStart compiles + runs via
# ts-node (dev) or pre-built JS (prod). This unit uses ts-node — same behavior
# as `npm run dev`. For prod, run `npm run build` first then point at dist/.

[Unit]
Description=Jovio API Server (Express + Razorpay + Wati)
After=network.target jovio-pipeline.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/jovi/api-server

# Secrets from: /etc/jovio/api-server.env (mode 0640, owner root:ubuntu)
# Contains: SUPABASE_*, RAZORPAY_*, WATI_*, INTERNAL_SECRET, PIPELINE_URL,
# RESEND_*, EXOTEL_WEBHOOK_TOKEN, RAZORPAY_PLAN_*
EnvironmentFile=/etc/jovio/api-server.env

# Run via ts-node for simplicity. For lower memory: precompile with
# `npm run build` then change ExecStart to: /usr/bin/node /home/ubuntu/jovi/api-server/dist/index.js
ExecStart=/home/ubuntu/jovi/api-server/node_modules/.bin/ts-node src/index.ts

Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

StandardOutput=append:/var/log/jovio-api.out.log
StandardError=append:/var/log/jovio-api.err.log

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/jovi/api-server /var/log

LimitNOFILE=65535
MemoryHigh=400M
MemoryMax=512M

[Install]
WantedBy=multi-user.target
