# infra/systemd/jovio-pipeline.service
# Jovio Voice Pipeline (Python FastAPI)
#
# Install:
#   sudo cp infra/systemd/jovio-pipeline.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now jovio-pipeline
#
# Replaces the previous Supervisor config at
# /etc/supervisor/conf.d/jovio-pipeline.conf — if migrating, first run:
#   sudo supervisorctl stop jovio-pipeline
#   sudo rm /etc/supervisor/conf.d/jovio-pipeline.conf
#   sudo supervisorctl reread && sudo supervisorctl update

[Unit]
Description=Jovio Voice Pipeline (FastAPI + LiveKit + Sarvam + Gemini)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/jovi/voice-pipeline

# Load secrets from a file outside the repo. Owner: root:ubuntu, mode 0640.
# Contains: SUPABASE_URL, SUPABASE_KEY, SARVAM_API_KEY, GEMINI_API_KEY,
# LIVEKIT_URL/KEY/SECRET, JOVIO_RECORDING_KEY, SUPABASE_RECORDINGS_BUCKET
EnvironmentFile=/etc/jovio/voice-pipeline.env

ExecStart=/home/ubuntu/jovi/voice-pipeline/venv/bin/uvicorn main:app \
          --host 0.0.0.0 --port 8000 --workers 2

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Logging — journald is captured by the CloudWatch agent via syslog OR
# we can write to files that the agent tails. Files are simpler.
StandardOutput=append:/var/log/jovio-pipeline.out.log
StandardError=append:/var/log/jovio-pipeline.err.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/jovi/voice-pipeline /var/log

# Resource limits
LimitNOFILE=65535
MemoryHigh=1500M
MemoryMax=1800M

[Install]
WantedBy=multi-user.target
