# infra/systemd/jovio-onboarding.timer
# Daily timer for the onboarding email job. Runs at 10:00 IST (04:30 UTC).
#
# That's the sweet spot in Indian B2B: late enough that recipients are
# at their desk, early enough that a "trial ending" email gives them
# a full business day to act before the cutoff.

[Unit]
Description=Daily trigger for Jovio onboarding emails

[Timer]
# 04:30 UTC = 10:00 IST
OnCalendar=*-*-* 04:30:00 UTC
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
